# fixes/fix-lid-actions/fix-lid-actions.service
[Unit]
Description=Fix unexpected lid actions

[Service]
ExecStart=/usr/bin/setkeycodes e057 245 e058 245

[Install]
WantedBy=multi-user.target

